# .local/scripts/udev/ac_plug_unplug.rules
ACTION=="change", SUBSYSTEM=="power_supply", KERNEL=="ADP1", RUN+="/bin/sh /home/ncho/.local/scripts/update-bar-battery"
